# packages/tempuscator/tempuscator-0.3.0a0-py3-none-any.whl/tempuscator/etc/systemd/system/obfuscator-watcher.service
[Unit]
Description=Obfuscator inotify daemon
After=network.target

[Service]
PIDFile=/run/obfuscator.pid
ExecStart=/usr/local/bin/watchfs.py

[Install]
WantedBy=multi-user.target

; vim:ft=systemd
